# v2/deploy/supervisord.conf
; /etc/supervisor/conf.d/rrp-v2.conf
;
; Supervisor configuration for the RRP v2 Docker stack.
; Supervisor keeps docker compose running and restarts it if it dies.
;
; Prerequisites:
;   sudo apt-get install -y supervisor
;   sudo cp rrp-v2.conf /etc/supervisor/conf.d/rrp-v2.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update
;
; Useful commands:
;   sudo supervisorctl status
;   sudo supervisorctl restart rrp-v2-stack
;   sudo supervisorctl tail -f rrp-v2-stack

; ── Docker Compose stack ──────────────────────────────────────────────────────
; Runs `docker compose up` in the foreground so supervisor can track the PID.
; Docker Compose itself manages app / worker / postgres / redis.

[program:rrp-v2-stack]
command=/usr/bin/docker compose -f /opt/rrp-v2/docker-compose.yml up --no-color
directory=/opt/rrp-v2
autostart=true
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
user=azureadmin
stdout_logfile=/var/log/supervisor/rrp-v2-stack.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/rrp-v2-stack.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

; ── Container health watchdog ─────────────────────────────────────────────────
; Polls the container health endpoint every 60 s and force-restarts the stack if
; the API stops responding (exits non-zero → supervisor triggers autorestart of
; rrp-v2-stack via `supervisorctl restart rrp-v2-stack`).
;
; This is intentionally a lightweight script, not a long-running process,
; so exitcodes-nodaemon and autorestart=unexpected keep it cycling.

[program:rrp-v2-watchdog]
command=/opt/rrp-v2/deploy/watchdog.sh
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=0
startretries=999
user=azureadmin
stdout_logfile=/var/log/supervisor/rrp-v2-watchdog.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=2
stderr_logfile=/var/log/supervisor/rrp-v2-watchdog.err

[group:rrp-v2]
programs=rrp-v2-stack,rrp-v2-watchdog
